# docker/crontab
* * * * * curl -X GET http://voicevox-timebell:8080/api/alert
